# conf/hotsearch_super.conf
[unix_http_server]
file=/tmp/supervisor.sock

[inet_http_server] 
port=127.0.0.1:9001
username=xiaomo
password=19940809

[supervisord]
logfile=/tmp/supervisord.log 
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid 
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[group:tornadoes]
programs=tornado-8010,tornado-8011,tornado-8012,tornado-8013

[program:tornado-8010]
# directory:程序的启动目录
directory=/spiders/weibohotsearch_crawler/hotsearch_api/
# command:启动命令(/spiders/env_py2.7/bin/python:python虚拟环境位置)
command=/spiders/env_py2.7/bin/python server.py --port=8010
# user:选择用户启动
user=root
# autorestart:程序异常退出后自动重启
autorestart=true
# killasgroup:进程收到 kill 信号时，会自动将该信号发给进程的子进程
stopasgroup=true
killasgroup=true
# 把 stderr 重定向到 stdout
redirect_stderr=true
# stdout 日志文件位置
stdout_logfile=/tmp/tornado_api_8010.log
loglevel=info

[program:tornado-8011]
directory=/spiders/weibohotsearch_crawler/hotsearch_api/
command=/spiders/env_py2.7/bin/python server.py --port=8011
user=root
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/tmp/tornado_api_8011.log
loglevel=info

[program:tornado-8012]
directory=/spiders/weibohotsearch_crawler/hotsearch_api/
command=/spiders/env_py2.7/bin/python server.py --port=8012
user=root
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/tmp/tornado_api_8012.log
loglevel=info

[program:tornado-8013]
directory=/spiders/weibohotsearch_crawler/hotsearch_api/
command=/spiders/env_py2.7/bin/python server.py --port=8013
user=root
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/tmp/tornado_api_8013.log
loglevel=info
